# ansible/roles/syslog/files/syslog-rotate.cron
30 1 * * * root /usr/local/sbin/syslog-rotate.sh

