# ansible/roles/tcpdumpd/files/tcpdump-rotate.cron
30 1 * * * root /usr/local/sbin/tcpdump-rotate.sh

