# ansible/roles/tcpdumpd/files/tcpdumpd-watcher.path
[Path]
PathModified=/etc/tcpdumpd/

[Install]
WantedBy=multi-user.target